# ai-assistant.service
[Unit]
Description=AI Assistant Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%%USER%%
WorkingDirectory=%%CONFIG_DIR%%
EnvironmentFile=%%CONFIG_DIR%%/.env
ExecStart=%%NODE_PATH%% %%CLI_PATH%% start
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-assistant

[Install]
WantedBy=multi-user.target
